# files/etc/init.d/mwan4
#!/bin/sh /etc/rc.common
# SPDX-License-Identifier: AGPL-3.0-or-later
# Copyright 2026 MOSSDeF, Stan Grishin (stangri@melmac.ca).
# Based on original mwan3 by Florian Eckert <fe@dev.tdt.de>

readonly PKG_VERSION='dev'

. "${IPKG_INSTROOT}/lib/functions.sh"

START=19
USE_PROCD=1

MWAN4_STATUS_DIR="/var/run/mwan4"

service_running() {
	[ -d "$MWAN4_STATUS_DIR" ]
}

start_tracker_family() {
	local interface="$1"
	local family="$2"

	procd_open_instance "track_${interface}_${family}"
	procd_set_param command ucode -S -L /lib/mwan4 /lib/mwan4/mwan4track.uc -- "$interface" "$family"
	procd_set_param respawn
	procd_close_instance
}

start_tracker() {
	local enabled track_ip families interface
	interface=$1
	config_get_bool enabled "$interface" 'enabled' '0'
	config_get track_ip "$interface" 'track_ip'
	[ "$enabled" -eq 0 ] && return
	[ -z "$track_ip" ] && return

	# Get family list (default: ipv4)
	config_get families "$interface" 'family'
	[ -z "$families" ] && families="ipv4"

	local family
	for family in $families; do
		start_tracker_family "$interface" "$family"
	done
}

start_service() {
	config_load mwan4
	config_foreach start_tracker interface

	# Delegate nftables/routes/rules setup to ucode
	ucode -S -L /lib/mwan4 /lib/mwan4/init_setup.uc

	procd_open_instance rtmon_ipv4
	procd_set_param command ucode -S -L /lib/mwan4 /lib/mwan4/mwan4rtmon.uc -- ipv4
	procd_set_param respawn
	procd_close_instance

	# Check IPv6 support
	if command -v ip > /dev/null && ip -6 addr show > /dev/null 2>&1; then
		procd_open_instance rtmon_ipv6
		procd_set_param command ucode -S -L /lib/mwan4 /lib/mwan4/mwan4rtmon.uc -- ipv6
		procd_set_param respawn
		procd_close_instance
	fi
}

stop_service() {
	# Delegate teardown to ucode
	ucode -S -L /lib/mwan4 /lib/mwan4/init_teardown.uc
}

service_stopped() { procd_set_config_changed firewall; }

reload_service() {
	stop
	start
}

service_triggers() {
	procd_add_reload_trigger 'mwan4'
}
